# supervisord.conf
[supervisord]
nodaemon=true                    ; start in foreground if true; default false

[program:dovecot]
command=/usr/sbin/dovecot -F              ; the program (relative uses PATH, can take args)

[program:rsyslog]
command=rsyslogd -n              ; the program (relative uses PATH, can take args)
